# family-piframe.cron
SHELL=/bin/bash
*/2 * * * * export DISPLAY=:0 && sudo /usr/bin/python3 /home/pi/zero-fun-net/family-piframe/main.py &> /home/pi/zero-fun-net/family-piframe/console.log &
